# rclone-bisync-gdrive.service
[Unit]
Description=Google Drive Sync Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rclone-bisync-gdrive.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=timers.target
